# systemd/speedtest.crontab
*/5 * * * * /etc/zabbix/bin/zbx-speedtest.sh
